# roles/Monitoring/monitoring_ansible/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=ansible-admin
Group=ansible-admin
Type=simple
ExecStart=/usr/local/bin/node_exporter

[Install]
WantedBy=multi-user.target
